# systemd/dctl-image-build.service
[Unit]
Description=Rebuild dctl devcontainer base images

[Service]
Type=oneshot
ExecStart=dctl image build --all
